# crontab

# This updates our crontab every minute automatically
* * * * * crontab /bigmetadata/crontab

